# distros/openwrt/clarissa.init
#!/bin/sh /etc/rc.common

# part of the clarissa project
# https://gitlab.com/evils/clarissa

START=32
STOP=16

USE_PROCD=1

start_service()
{
	procd_open_instance
	procd_set_param command /opt/clarissa/clarissa/clarissa_static
	procd_append_param command -pI br-lan
	procd_set_param respawn
	procd_close_instance
}
